# cronie/runtime/etc/periodic/hourly/task
#!/bin/sh
echo $(date) -hourly- >> /cronie.log