# cronie/runtime/etc/periodic/15min/task
#!/bin/sh
echo $(date) -15min- >> /cronie.log