# cronie/runtime/etc/periodic/monthly/task
#!/bin/sh
echo $(date) -montly- >> /cronie.log