# cronie/runtime/etc/periodic/daily/task
#!/bin/sh
echo $(date) -daily- >> /cronie.log