# cronie/runtime/etc/periodic/weekly/task
#!/bin/sh
echo $(date) -weekly- >> /cronie.log
echo -----rolllback----- > /cronie.log